# deploy/laravel-reverb.service
[Unit]
Description=Laravel Reverb WebSocket Server
After=network.target redis.service

[Service]
User=www-data
Group=www-data
Restart=always
ExecStart=/usr/bin/php8.4 /var/www/msw/artisan reverb:start --host=0.0.0.0 --port=8080
StandardOutput=append:/var/log/laravel-reverb.log
StandardError=append:/var/log/laravel-reverb.log

[Install]
WantedBy=multi-user.target
